# config/frpc.supervisor.ws.banner.conf
[program:frpc_ws_argue]
command=/root/center_frpc_ws/bin/frpc argue 50001
stopasgroup=true
killasgroup=true
restart=always
user=root
