# config/frpc.supervisor.ws.argue.conf
[program:argue_ws]
command=/root/center_frpc_ws/bin/frpc ws_argue 4321
stopasgroup=true
killasgroup=true
restart=always
user=root
