# config/frpc.supervisor.ws.neo4j.conf
[program:frpc_ws_neo4j]
command=/root/center_frpc_ws/bin/frpc neo4j 50001
stopasgroup=true
killasgroup=true
restart=always
user=root
